# ansible/roles/op-stack/templates/lib/systemd/system/op-node.service.j2
[Unit]
Description=Op Node
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=360
Restart=always
ExecStartPre=-/usr/bin/docker stop --time=300 op-node
ExecStartPre=-/usr/bin/docker rm op-node
ExecStartPre=/usr/bin/docker pull us-docker.pkg.dev/oplabs-tools-artifacts/images/op-node:{{- op_stack_op_node_version }}
ExecStart=/usr/bin/docker run \
    --name op-node \
    --rm \
    --network host \
    -v /home/op-stack/jwt.txt:/jwt.txt \
    -v /home/op-stack/rollup.json:/rollup.json \
    --entrypoint op-node \
    us-docker.pkg.dev/oplabs-tools-artifacts/images/op-node:{{- op_stack_op_node_version }} \
        --l2=http://localhost:8551 \
        --l2.jwt-secret=/jwt.txt \
        --sequencer.enabled \
        --sequencer.l1-confs=5 \
        --verifier.l1-confs=4 \
        --rollup.config=/rollup.json \
        --rpc.addr=0.0.0.0 \
        --rpc.port=8547 \
        --p2p.disable \
        --rpc.enable-admin \
        --p2p.sequencer.key={{- op_stack_GS_SEQUENCER_PRIVATE_KEY }} \
        --l1={{- op_stack_L1_RPC_URL }} \
        --l1.beacon={{- op_stack_L1_BEACON_URL }} \
        --l1.rpckind={{- op_stack_L1_RPC_KIND }}


[Install]
WantedBy=multi-user.target